# launchd/com.kai.voice-inbox.poll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.kai.voice-inbox.poll</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>cd __PROJECT_DIR__ &amp;&amp; ./dist/voice-inbox poll --once --json &gt;&gt; "$HOME/Library/Logs/voice-inbox-daemon/poll.log" 2&gt;&gt; "$HOME/Library/Logs/voice-inbox-daemon/poll.err.log"</string>
    </array>
  </dict>
</plist>
